# etc/udev/rules.d/30_rockchip_ebc.rules
# old driver
# KERNELS=="rockchip_ebc", RUN+="/bin/chgrp video /sys/module/rockchip_ebc/parameters/force_refresh", RUN+="/bin/chmod g+w /sys/module/rockchip_ebc/parameters/force_refresh", RUN+="/bin/chgrp video /sys/module/rockchip_ebc/parameters/lut_type", RUN+="/bin/chmod g+w /sys/module/rockchip_ebc/parameters/lut_type"

# default_waveform  diff_mode  direct_mode  panel_reflection  skip_reset
#KERNELS=="rockchip_ebc", RUN+="/bin/chgrp video /sys/module/rockchip_ebc/parameters/default_waveform", RUN+="/bin/chmod g+w /sys/module/rockchip_ebc/parameters/default_waveform", RUN+="/bin/chgrp video /sys/module/rockchip_ebc/parameters/diff_mode", RUN+="/bin/chmod g+w /sys/module/rockchip_ebc/parameters/diff_mode"

SUBSYSTEM=="module", KERNEL=="rockchip_ebc", RUN+="/bin/chgrp video /sys/module/rockchip_ebc/parameters/auto_refresh /sys/module/rockchip_ebc/parameters/bw_dither_invert /sys/module/rockchip_ebc/parameters/bw_threshold /sys/module/rockchip_ebc/parameters/bw_mode /sys/module/rockchip_ebc/parameters/diff_mode /sys/module/rockchip_ebc/parameters/split_area_limit /sys/module/rockchip_ebc/parameters/refresh_waveform /sys/module/rockchip_ebc/parameters/refresh_threshold /sys/module/rockchip_ebc/parameters/default_waveform /sys/module/rockchip_ebc/parameters/limit_fb_blits /sys/module/rockchip_ebc/parameters/prepare_prev_before_a2 /sys/module/rockchip_ebc/parameters/skip_reset /sys/module/rockchip_ebc/parameters/direct_mode", RUN+="/bin/chmod g+w /sys/module/rockchip_ebc/parameters/direct_mode /sys/module/rockchip_ebc/parameters/skip_reset /sys/module/rockchip_ebc/parameters/prepare_prev_before_a2 /sys/module/rockchip_ebc/parameters/limit_fb_blits /sys/module/rockchip_ebc/parameters/bw_dither_invert /sys/module/rockchip_ebc/parameters/bw_threshold /sys/module/rockchip_ebc/parameters/bw_mode /sys/module/rockchip_ebc/parameters/default_waveform /sys/module/rockchip_ebc/parameters/diff_mode /sys/module/rockchip_ebc/parameters/auto_refresh /sys/module/rockchip_ebc/parameters/split_area_limit /sys/module/rockchip_ebc/parameters/refresh_waveform /sys/module/rockchip_ebc/parameters/refresh_threshold"
DRIVER=="rockchip-ebc", RUN+="/bin/chgrp video /sys/%p/power/control", RUN+="/bin/chmod g+w /sys/%p/power/control"
